# roles/gpsclock/templates/gpsd.service
[Service]
ExecStartPre=/usr/bin/stty -F /dev/{{ serialport }} speed {{ serialbaud }}
